# crontab
# 29/08/2025 às 13:45 — executa uma ÚNICA vez
45 13 29 8 * root  cd /app && if [ -d .git ]; then GIT_TERMINAL_PROMPT=0 git pull --rebase --autostash || true; fi; /usr/local/bin/python /app/main.py --all >> /var/log/cron.log 2>&1
